# uniguard-bridge.service
# ─────────────────────────────────────────────────────────────────────────────
# NOTE: install.sh generates the actual systemd unit at install time, templated
# with the correct service user and APP_DIR. This file is a reference copy.
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Uniguard Pro Bridge — RTSP to HLS Streaming Gateway
Documentation=https://github.com/uniguardpro-monitoring/uniguard-pro-bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/uniguard-bridge
ExecStart=/opt/uniguard-bridge/venv/bin/uvicorn app.main:app \
          --host 0.0.0.0 \
          --port 8080 \
          --workers 1 \
          --log-level info

Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60

User=uniguard
Group=uniguard

NoNewPrivileges=true
PrivateTmp=true

Environment="UGBRIDGE_HOST=0.0.0.0"
Environment="UGBRIDGE_PORT=8080"
Environment="UGBRIDGE_STREAM_TIMEOUT_SECONDS=300"

[Install]
WantedBy=multi-user.target
