# uniguard-bridge-update.timer
# Reference copy — the actual unit is generated by install.sh
[Unit]
Description=Uniguard Pro Bridge — Auto-Update Timer

[Timer]
OnCalendar=*-*-* 11:00:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
